# SRC/DEB_picopilot-scripts/etc/supervisor/conf.d/nicPrep.conf
[program:nicPrep]
autostart=false
command=/bin/bash /opt/piCopilot-scripts/nicPrep.sh
directory=/opt/piCopilot-scripts
stderr_logfile=/tmp/nicERR
stdout_logfile=/tmp/nicOUT
