# SRC/DEB_picopilot-unmanned/etc/supervisor/conf.d/copilot.conf
[program:copilot]
autostart=true
command=/opt/piCopilot-unmanned/piController.py
directory=/opt/piCopilot-unmanned
stderr_logfile=/tmp/unmannedERR
stdout_logfile=/tmp/unmannedOUT
