# SRC/DEB_picopilot-unmanned/etc/supervisor/conf.d/telemetry_Service.conf
[program:telemetry_Service]
autostart=false
command=bash ./mproxy.wrapper
directory=/opt/piCopilot-unmanned/modules/telemetry
stderr_logfile=/tmp/telemSvcERR
stdout_logfile=/tmp/telemSvcOUT
