# SRC/DEB_picopilot-unmanned/DEBIAN/postinst
#!/usr/bin/env bash

## Start supervisor daemon
systemctl start supervisor
systemctl enable supervisor
