# SRC/DEB_picopilot-idrop/etc/supervisor/conf.d/idrop.conf
[program:idropController]
autostart=true
command=/usr/bin/python3 /opt/piCopilot-idrop/idrop.py
directory=/opt/piCopilot-idrop
stderr_logfile=/tmp/idropERR
stdout_logfile=/tmp/idropOUT
