# SRC/DEB_picopilot-idrop/etc/supervisor/conf.d/kSnarfPsql.conf
[program:kSnarfPsql]
autostart=false
command=/usr/bin/python3 /opt/piCopilot-idrop/kSnarf.py
directory=/opt/piCopilot-idrop
stderr_logfile=/tmp/kSnarfERR
stdout_logfile=/tmp/kSnarfOUT
