# SRC/DEB_picopilot-idrop/etc/supervisor/conf.d/nicMon.conf
[program:nicMon]
autostart=false
command=airmon-ng start wlan1
directory=/opt/piCopilot-idrop
stderr_logfile=/tmp/monERR
stdout_logfile=/tmp/monOUT
